# config/supervisor/3-vnc.conf
[program:VNC]
command=bash -c 'x11vnc -rfbport 5900 -ncache 10 -geometry $(jq --raw-output ".resolution // \"768x1024\"" /data/options.json) -forever -shared -alwaysshared -bg -desktop "$VNC_TITLE"'
startsecs=0
autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=3
